# docker/runit/qmailapi/run
#!/bin/sh
exec python3 /opt/qmailapi/app.py
