# docker/runit/qmail-smtpd/run
#!/bin/sh
QMAILDIR=$(getent passwd qmaild | cut -d: -f6)
QMAILDUID=$(id -u vpopmail)
NOFILESGID=$(id -g vpopmail)
MAXSMTPD=$(cat $QMAILDIR/control/concurrencyincoming)
SOFTLIMIT=$(cat $QMAILDIR/control/softlimit)
LOCAL=$(head -1 $QMAILDIR/control/me)
TCPRULES_DIR=$QMAILDIR/control

export CHKUSER_START=ALWAYS
export FORCEAUTHMAILFROM=1

# Greetdelay — configurable via control/greetdelay (default: 5)
export SMTPD_GREETDELAY=$(cat $QMAILDIR/control/greetdelay 2>/dev/null || echo 5)
export DROP_PRE_GREET=1

# SURBL — configurable via control/surbl (default: 0 = disabled)
export SURBL=$(cat $QMAILDIR/control/surbl 2>/dev/null || echo 0)

# DKIM verification for inbound mail; simscan (AV + spam) runs after
export QMAILQUEUE=$QMAILDIR/bin/qmail-dkim
export DKIMVERIFY="FGHKLMNOQRTVWp"
export RELAYCLIENT_NODKIMVERIFY=1
if [ "${SIMSCAN_ENABLE:-true}" = "true" ]; then
    export DKIMQUEUE=$QMAILDIR/bin/simscan
    export SIMSCAN_DEBUG=${SIMSCAN_DEBUG:-0}
fi

if [ ! -f $QMAILDIR/control/smtpplugins ] || [ ! -d $QMAILDIR/plugins ]; then
    export ENABLE_SPP=0
else
    export ENABLE_SPP=1
fi

DUALSTACK=$(cat $QMAILDIR/control/dualstack 2>/dev/null || echo 0)
[ "$DUALSTACK" = "1" ] && BINDADDR="::" || BINDADDR="0"

# jgreylist — file-based binary wrapper (control/jgreylist=1 to enable)
JGREYLIST=$(cat $QMAILDIR/control/jgreylist 2>/dev/null || echo 0)
export JGREYLIST_DIR="$QMAILDIR/jgreylist"
export JGREYLIST_LOG_SMTP=1

# qmail-spp greylisting plugin — MySQL-backed, runs inside smtpd via SPP.
# Activated by the presence of control/greylisting (its config file).
# The plugin checks for GREYLISTING in the environment; empty string = enabled.
if [ -f "$QMAILDIR/control/greylisting" ]; then
    export GREYLISTING=""
fi

if [ "$JGREYLIST" = "1" ]; then
    exec chpst -m "$SOFTLIMIT" \
        /usr/local/bin/tcpserver -v -R -l "$LOCAL" \
        -x ${TCPRULES_DIR}/tcp.smtp.cdb -c "$MAXSMTPD" \
        -u "$QMAILDUID" -g "$NOFILESGID" "$BINDADDR" 25 \
        $QMAILDIR/bin/jgreylist \
        $QMAILDIR/bin/qmail-smtpd \
        ~vpopmail/bin/vchkpw /bin/true 2>&1
else
    exec chpst -m "$SOFTLIMIT" \
        /usr/local/bin/tcpserver -v -R -l "$LOCAL" \
        -x ${TCPRULES_DIR}/tcp.smtp.cdb -c "$MAXSMTPD" \
        -u "$QMAILDUID" -g "$NOFILESGID" "$BINDADDR" 25 \
        $QMAILDIR/bin/qmail-smtpd \
        ~vpopmail/bin/vchkpw /bin/true 2>&1
fi
