# docker/runit/qmail-smtpd/run
#!/bin/sh
QMAILDIR=$(getent passwd qmaild | cut -d: -f6)
QMAILDUID=$(id -u vpopmail)
NOFILESGID=$(id -g vpopmail)
MAXSMTPD=$(cat $QMAILDIR/control/concurrencyincoming)
SOFTLIMIT=$(cat $QMAILDIR/control/softlimit)
LOCAL=$(head -1 $QMAILDIR/control/me)
TCPRULES_DIR=$QMAILDIR/control

export CHKUSER_START=ALWAYS
export FORCEAUTHMAILFROM=1

# Greetdelay — configurable via control/greetdelay (default: 5)
export SMTPD_GREETDELAY=$(cat $QMAILDIR/control/greetdelay 2>/dev/null || echo 5)
export DROP_PRE_GREET=1

# SURBL — configurable via control/surbl (default: 0 = disabled).
# When enabled, surblqueue is prepended to the queue chain so it can inspect message
# bodies for SURBL-listed URIs before DKIM signing / scanning.
SURBL_VAL=$(cat $QMAILDIR/control/surbl 2>/dev/null || echo 0)
if [ "$SURBL_VAL" = "1" ]; then
    export SURBL=1
    export QMAILQUEUE=$QMAILDIR/bin/surblqueue
    export SURBLQUEUE=$QMAILDIR/bin/qmail-dkim
else
    export QMAILQUEUE=$QMAILDIR/bin/qmail-dkim
fi

# DKIM verification for inbound mail; simscan (AV + spam) runs after
# DKIMVERIFY read from control/dkimverify — written by entrypoint based on DKIM_VERIFY_LAYER.
# Empty = disabled (rspamd layer); "FGHKLMNOQRTVWp" = full verification (qmail layer).
export DKIMVERIFY=$(cat $QMAILDIR/control/dkimverify 2>/dev/null || true)
export RELAYCLIENT_NODKIMVERIFY=1
if [ "${SIMSCAN_ENABLE:-true}" = "true" ]; then
    export DKIMQUEUE=$QMAILDIR/bin/simscan
    export SIMSCAN_DEBUG=${SIMSCAN_DEBUG:-0}
fi

if [ -f $QMAILDIR/control/smtpplugins ] && [ -d $QMAILDIR/plugins ]; then
    export ENABLE_SPP=1
fi

DUALSTACK=$(cat $QMAILDIR/control/dualstack 2>/dev/null || echo 0)
[ "$DUALSTACK" = "1" ] && BINDADDR="::" || BINDADDR="0"

# jgreylist — file-based binary wrapper (control/jgreylist=1 to enable)
JGREYLIST=$(cat $QMAILDIR/control/jgreylist 2>/dev/null || echo 0)
export JGREYLIST_DIR="$QMAILDIR/jgreylist"
export JGREYLIST_LOG_SMTP=1

# qmail-spp greylisting plugin — MySQL-backed, runs inside smtpd via SPP.
# GREYLISTING="": activates the plugin (required by greylisting.c).
# GLCONFIGFILE: path to the greylisting config file (relative to qmail dir,
# since qmail-smtpd chdirs to /var/qmail at startup).
# Note: the plugin automatically skips RELAYCLIENT connections.
if [ -f "$QMAILDIR/control/greylisting" ]; then
    export GREYLISTING=""
    export GLCONFIGFILE=control/greylisting
fi

if [ "$JGREYLIST" = "1" ]; then
    exec chpst -m "$SOFTLIMIT" \
        /usr/local/bin/tcpserver -v -R -l "$LOCAL" \
        -x ${TCPRULES_DIR}/tcp.smtp.cdb -c "$MAXSMTPD" \
        -u "$QMAILDUID" -g "$NOFILESGID" "$BINDADDR" 25 \
        $QMAILDIR/bin/jgreylist \
        $QMAILDIR/bin/qmail-smtpd \
        ~vpopmail/bin/vchkpw /bin/true 2>&1
else
    exec chpst -m "$SOFTLIMIT" \
        /usr/local/bin/tcpserver -v -R -l "$LOCAL" \
        -x ${TCPRULES_DIR}/tcp.smtp.cdb -c "$MAXSMTPD" \
        -u "$QMAILDUID" -g "$NOFILESGID" "$BINDADDR" 25 \
        $QMAILDIR/bin/qmail-smtpd \
        ~vpopmail/bin/vchkpw /bin/true 2>&1
fi
