# docker/runit/qmail-submission/run
#!/bin/sh
QMAILDIR=$(getent passwd qmaild | cut -d: -f6)
QMAILDUID=$(id -u vpopmail)
NOFILESGID=$(id -g vpopmail)
MAXSMTPD=$(cat $QMAILDIR/control/concurrencyincoming)
SOFTLIMIT=$(cat $QMAILDIR/control/softlimit)
LOCAL=$(head -1 $QMAILDIR/control/me)
TCPRULES_DIR=$QMAILDIR/control
DUALSTACK=$(cat $QMAILDIR/control/dualstack 2>/dev/null || echo 0)
[ "$DUALSTACK" = "1" ] && BINDADDR="::" || BINDADDR="0"

# Auth required on port 587
export SMTPAUTH="!"
export FORCETLS=1
export CHKUSER_START=ALWAYS
export FORCEAUTHMAILFROM=1
export SMTPD_GREETDELAY=$(cat $QMAILDIR/control/greetdelay 2>/dev/null || echo 5)
export DROP_PRE_GREET=1

# SURBL and DKIM — same as port 25
export SURBL=$(cat $QMAILDIR/control/surbl 2>/dev/null || echo 0)
export QMAILQUEUE=$QMAILDIR/bin/qmail-dkim
if [ "${SIMSCAN_ENABLE:-true}" = "true" ]; then
    export DKIMQUEUE=$QMAILDIR/bin/simscan
    export SIMSCAN_DEBUG=${SIMSCAN_DEBUG:-0}
fi

# SPP plugins (e.g. greylisting) — only set when plugins are actually available.
# Setting ENABLE_SPP=0 still triggers the SPP code path in qmail-smtpd; leave it unset to disable.
if [ -f $QMAILDIR/control/smtpplugins ] && [ -d $QMAILDIR/plugins ]; then
    export ENABLE_SPP=1
fi

# Rate limiting — restrict how many messages an auth user/domain/IP may send.
# control/relaylimits format: "user@domain:N", "domain:N", "IP:N", ":N" (default)
# 0 = unlimited for that entry. Counts reset by cron via greylisting_cleanup or similar.
export RCPTCHECK="$QMAILDIR/bin/rcptcheck-overlimit"
export RCPTCHECKRELAYCLIENT="1"
export OVERLIMITDIR="$QMAILDIR/overlimit"
export LIMITSCONTROLFILE="$QMAILDIR/control/relaylimits"

exec chpst -m "$SOFTLIMIT" \
    /usr/local/bin/tcpserver -v -R -l "$LOCAL" \
    -x ${TCPRULES_DIR}/tcp.submission.cdb -c "$MAXSMTPD" \
    -u "$QMAILDUID" -g "$NOFILESGID" "$BINDADDR" 587 \
    $QMAILDIR/bin/qmail-smtpd \
    ~vpopmail/bin/vchkpw /bin/true 2>&1
